# data/launchd/com.jobstack.scheduler.live.pipeline.planning.only.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AbandonProcessGroup</key>
	<true/>
	<key>Label</key>
	<string>com.jobstack.scheduler.live.pipeline.planning.only</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/sriramhariharanneelakantan/.venvs/job_scrap314/bin/python</string>
		<string>-u</string>
		<string>-m</string>
		<string>src.pipeline.scheduler</string>
		<string>--job</string>
		<string>live_pipeline</string>
		<string>--history-path</string>
		<string>data/scheduler_run_history.jsonl</string>
		<string>--planning-only</string>
		<string>--output-dir</string>
		<string>outputs/application_planning</string>
		<string>--job-limit</string>
		<string>50</string>
		<string>--job-packet-limit</string>
		<string>0</string>
		<string>--llm-actions</string>
		<string>APPLY,APPLY_REVIEW_VARIANTS</string>
		<string>--delete-seen-data</string>
		<string>no</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>outputs/scheduler_logs/com.jobstack.scheduler.live.pipeline.planning.only.err.log</string>
	<key>StandardOutPath</key>
	<string>outputs/scheduler_logs/com.jobstack.scheduler.live.pipeline.planning.only.out.log</string>
	<key>StartInterval</key>
	<integer>14400</integer>
	<key>WorkingDirectory</key>
	<string>/Users/sriramhariharanneelakantan/Desktop/Gen AI learning/Job_board_scrapping</string>
</dict>
</plist>
