# data/launchd/com.jobstack.scheduler.agent.discovery.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>AbandonProcessGroup</key>
	<true/>
	<key>Label</key>
	<string>com.jobstack.scheduler.agent.discovery</string>
	<key>ProcessType</key>
	<string>Background</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/sriramhariharanneelakantan/.venvs/job_scrap314/bin/python</string>
		<string>-u</string>
		<string>-m</string>
		<string>src.pipeline.scheduler</string>
		<string>--job</string>
		<string>agent_discovery</string>
		<string>--history-path</string>
		<string>data/scheduler_run_history.jsonl</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>outputs/scheduler_logs/com.jobstack.scheduler.agent.discovery.err.log</string>
	<key>StandardOutPath</key>
	<string>outputs/scheduler_logs/com.jobstack.scheduler.agent.discovery.out.log</string>
	<key>StartInterval</key>
	<integer>21600</integer>
	<key>WorkingDirectory</key>
	<string>/Users/sriramhariharanneelakantan/Desktop/Gen AI learning/Job_board_scrapping</string>
</dict>
</plist>
